# .config/systemd/user-avail/docker-gc.service
[Unit]
Description=Docker garbage collection

[Service]
ExecStart=/home/trevorj/go/bin/docker-gc

[Install]
WantedBy=default.target
